# src/deb/control/prerm
#!/bin/sh
service sspc stop
update-rc.d -f sspc remove
